# contrib/upstart/reserver.conf
description "Crusher Worker"

setuid crusher
start on filesystem or runlevel [2345]
stop on shutdown

script
    source venv/bin/activate
    python app.py
end script
